# experiments/internetscanner/templates/outlook.service.j2
[Unit]
Description=sticky elephant run script
Wants=network-online.target
After=network-online.target

[Service]
Restart=always
RestartSec=1s

WorkingDirectory=/home/{{ ansible_user_id }}/owa-honeypot
ExecStart=/bin/bash -c 'python3 owa_pot.py'

TimeoutStartSec=300
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
